# aws/hdf/resources/scripts/nifi-connectors.path
[Unit] 
Description=NiFi-Home-Listener
[Path]
PathModified=/home/nifi/